# contrib/com.syfrah.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syfrah.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/syfrah</string>
        <string>start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/syfrah.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/syfrah.log</string>

    <key>WorkingDirectory</key>
    <string>/tmp</string>
</dict>
</plist>
